# nginx/antmedia-cluster-check.service
[Unit]
Description=Ant Media Cluster Check Script

[Service]
ExecStart=/usr/bin/nginx-auto-scale.sh

[Install]
WantedBy=multi-user.target
